# glorytun-udp@.service.in
[Unit]
Description=Glorytun UDP on %I
After=network.target network-online.target

[Service]
Type=simple
Restart=always
EnvironmentFile=/etc/glorytun-udp/%i
ExecStart=/usr/local/bin/glorytun-udp-run /etc/glorytun-udp/%i
ExecStartPost=-/etc/glorytun-udp/post.sh /etc/glorytun-udp/%i
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
