# glorytun-tcp@.service.in
[Unit]
Description=Glorytun TCP on %I
After=network.target network-online.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/glorytun-tcp-run /etc/glorytun-tcp/%i
ExecStartPost=/usr/local/bin/omr-6in4 /etc/glorytun-tcp/%i start
ExecStopPost=/usr/local/bin/omr-6in4 /etc/glorytun-tcp/%i stop
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
