# mlvpn@.service.in
[Unit]
Description=MLVPN connection to %i
PartOf=mlvpn.service
ReloadPropagatedFrom=mlvpn.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/sbin/mlvpn --config /etc/mlvpn/%i.conf --name %i --user mlvpn --quiet
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/etc/mlvpn
Restart=always

[Install]
WantedBy=multi-user.target


