# omr-bypass.timer.in
[Unit]
Description=Timer for omr-bypass

[Timer]
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
